# init.d/statspoller
#!/bin/bash
# chkconfig: 35 20 80
# description: simple daemon wrapper to manage start/stopping of statspoller

export process_name="StatsPoller"
source /opt/$process_name/conf/statspoller.profile

NAME=`basename $0`
PID=`pgrep -f StatsPoller.jar`

case "$1" in
    start)
          cd /opt/StatsPoller
          if [ -z "$PID" ]; then
              (java $JAVA_OPTS > /dev/null 2>&1 &)
          else
              printf "Already running on PID $PID!\n"
              exit 1
          fi
          sleep 1
          ;;
          status)
              if [ -n "$PID" ]; then
                  ps -fp $PID
              else
                  printf "Serivice isn't running\n"
                  exit 1
              fi
              ;;
          stop)
              if [ -n "$PID" ]; then
                  kill $PID
              else
                  printf "Service isn't running\n"
                  exit 1
              fi
              ;;
          restart)
              $0 stop
              sleep 5
              $0 start
              ;;
          *)
              printf "usage: $0 {status|start|stop|restart}\n"
              exit 1
          esac
